# cron.d/cron
#!/bin/bash -l
00 8 * * * export $(cat /root/root.env | xargs); /go/src/app/main
