# cron.d/cron
#!/bin/bash -l
* * * * * export $(cat /root/root.env | xargs); /go/src/app/main
